# diagrams/conf.d/diagrams.conf
[program:diagrams]
command=%(ENV_HOME)s/projects/lyra-stack/diagrams/scripts/run.sh
directory=%(ENV_HOME)s/.agent/diagrams
autostart=true
autorestart=true
startsecs=3
startretries=3
stopwaitsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=%(ENV_HOME)s/projects/lyra-stack/diagrams/logs/diagrams.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=%(ENV_HOME)s/projects/lyra-stack/diagrams/logs/diagrams_error.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
